# bot.conf
[program:Ex_Minus_ut_Magis]
directory=/home/ubuntu/bot4
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/bot4/logfile_err.log
stdout_logfile=/home/ubuntu/bot4/logfile.log